# scripts/cloud-image/systemd/semiclaw-firstboot.service
[Unit]
Description=SemiClaw first-boot initialization (runs once)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target
Before=semiclaw.service
ConditionPathExists=/usr/local/sbin/semiclaw-firstboot.sh
ConditionPathExists=!/opt/SemiClaw/.firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/semiclaw-firstboot.sh
RemainAfterExit=no
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
